# files/release/box-air-release-init.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

@reboot           root    /usr/local/sbin/box-air-release-init
*/15 *    * * *   root    /usr/local/sbin/box-air-release-init
